# gemini-skill-browser-builder/image/supervisor/gemini-skill.conf
[program:gemini-skill-daemon]
directory=/opt/gemini-skill
command=/usr/bin/node /opt/gemini-skill/src/daemon/server.js
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/var/log/gemini-skill-daemon.log
stderr_logfile=/var/log/gemini-skill-daemon.err.log
environment=NODE_ENV="production",BROWSER_DEBUG_PORT="%(ENV_BROWSER_DEBUG_PORT)s",BROWSER_HEADLESS="%(ENV_BROWSER_HEADLESS)s",BROWSER_USER_DATA_DIR="%(ENV_BROWSER_USER_DATA_DIR)s",BROWSER_PROTOCOL_TIMEOUT="%(ENV_BROWSER_PROTOCOL_TIMEOUT)s",DAEMON_PORT="%(ENV_DAEMON_PORT)s",DAEMON_TTL_MS="%(ENV_DAEMON_TTL_MS)s",OUTPUT_DIR="%(ENV_OUTPUT_DIR)s"
